# os/src/config/overlay.d/50substrateos/usr/lib/systemd/system/ensure-oob-imagestore.service
[Unit]
ConditionPathExists=!/run/media/oob/imagestore

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=mkdir -m 0755 -p /run/media/oob/imagestore
