# images/chromestage/service/pulseaudio/run
#!/bin/sh
exec pulseaudio --verbose --exit-idle-time=-1 --disallow-exit
