# images/chromestage/service/chromium/run
#!/bin/sh
exec chromium \
  --disable-sandbox \
  --no-default-browser-check \
  --remote-allow-origins=* \
  --remote-debugging-port=9222 \
  --window-position=0,0 \
  --window-size=$CHROMIUM_WINDOW_SIZE \
  --no-first-run \
  --kiosk \
  --autoplay-policy=no-user-gesture-required \
  --disable-background-networking=true \
  --enable-features="NetworkService,NetworkServiceInProcess" \
  --disable-background-timer-throttling=true \
  --disable-backgrounding-occluded-windows=true \
  --disable-breakpad=true \
  --disable-client-side-phishing-detection=true \
  --disable-default-apps=true \
  --disable-dev-shm-usage=true \
  --disable-extensions=true \
  --disable-features="site-per-process,Translate,BlinkGenPropertyTrees" \
  --disable-hang-monitor=true \
  --disable-ipc-flooding-protection=true \
  --disable-popup-blocking=true \
  --disable-prompt-on-repost=true \
  --disable-renderer-backgrounding=true \
  --disable-sync=true \
  --force-color-profile="srgb" \
  --metrics-recording-only=true \
  --safebrowsing-disable-auto-update=true \
  --password-store="basic" \
  --use-mock-keychain=true \
  --disable-extensions \
  --disable-component-update \
  --disable-default-apps \
  --disable-hang-monitor \
  --disable-prompt-on-repost \
  --full-memory-crash-report \
  --log-level=0 \
  --enable-logging \
  --v=0 \
  $CHROMIUM_START_URL
# --start-maximized --start-fullscreen
