# os/fcos/usr/lib/systemd/system/nvidia-ctk-cdi-generate.service
[Unit]
ConditionPathExists=!/etc/cdi/nvidia.yaml

[Service]
Type=oneshot
ExecStart=nvidia-ctk cdi generate --output=/etc/cdi/nvidia.yaml
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
